# etc/s6-overlay/s6-rc.d/wordpress-seed/up
/etc/s6-overlay/scripts/wordpress-seed-up
