# etc/s6-overlay/s6-rc.d/mariadb-prepare/up
/etc/s6-overlay/scripts/mariadb-prepare
